# .config/profile.d/less.sh
#!/bin/sh
#
# zlib license
# Copyright (c) 2017 Victorien Elvinger
#

export LESS='--ignore-case -g -M -R -w --window=-5'
    # Default Less options
    # g: only highlight the last search (perf)
    # M: be more verbose than more (long prompt)
    # R: output ANSI "color" escape sequences in "raw" form
    # w: highlights  the  first  "new"  line after a forward

