# .config/profile.d/programs.sh
#!/bin/sh

export PAGER='less'
export VISUAL='vim'

